# files/rc.local
#!/bin/bash -e

cd /var/www/app/sns/

rake db:create RAILS_ENV=production
sqlite3 db/production.sqlite3 < /root/bad_sns_production.sql 

mkdir -p tmp/cache
chown -R www-data.www-data .
bundle exec unicorn_rails -c config/unicorn.rb -E production -D

/root/MailHog_linux_amd64 &

systemctl daemon-reload
